# resources/whois42d-ng.service
[Unit]
Description=Whois dn42 daemon

[Service]
ExecStart=/usr/local/bin/whois42d-ng --registry /var/lib/whois42d/registry
User=nobody
